# streamserver.service
[Unit]
Description=Raspberry Pi RTSP Stream Server with H.264 Hardware Encoding
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/streamserver
ExecStart=/usr/bin/python3 /home/pi/streamserver/streamserver.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits (increased for H.264 encoding)
LimitNOFILE=2048

# Camera and video access permissions
SupplementaryGroups=video

# Note: If using standard RTSP port 554, consider running with CAP_NET_BIND_SERVICE
# or change port to 8554 in config.py for non-root operation

[Install]
WantedBy=multi-user.target